# scripts/exia-robot.service
[Unit]
Description=Exia Ground Robot Autonomous + Radio Bridge
After=network.target local-fs.target
Wants=network.target

[Service]
Type=simple
User=exialabsargus
Group=exialabsargus
WorkingDirectory=/home/exialabsargus/exia_ws
ExecStart=/home/exialabsargus/exia_ws/scripts/exia_startup.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

Environment="HOME=/home/exialabsargus"
Environment="PATH=/home/exialabsargus/bin:/usr/local/bin:/usr/bin:/bin:/home/exialabsargus/.local/bin"

[Install]
WantedBy=multi-user.target
